# systemd_example/sdsb.timer
[Unit]
Description=Run SDSB daily

[Timer]
OnCalendar=*-*-* 09:00:00

